# scripts/bot-ai.service
[Unit]
Description=Bot AI Daemon
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/bot-ai
EnvironmentFile=-/opt/bot-ai/.env
ExecStart=/opt/bot-ai/.venv/bin/python /opt/bot-ai/src/daemon_main.py
Restart=always
RestartSec=2
User=root
Group=root

[Install]
WantedBy=multi-user.target
